# s6-rc.d/start-daemon/run
#!/command/execlineb -P
with-contenv
s6-setuidgid electrum
foreground { echo "Starting Electrum daemon in foreground mode..." }
/usr/local/bin/start-daemon.sh
