# s6-rc.d/service/run
#!/command/execlineb -P
with-contenv
s6-setuidgid electrum
foreground { echo "Starting Electrum service..." }
/usr/local/bin/service.sh
